# stage3_homebridge/02-wifi-connect/files/wifi-connect.service
[Unit]
Description=WiFi Captive Portal
After=NetworkManager.service
Before=nginx.service homebridge.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wifi-connect

[Install]
WantedBy=multi-user.target
