# stage3_custom/01-homebridge/files/hb-check.service
[Unit]
Description=hb-check service
Before=homebridge.service
After=nginx.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/hb-check

[Install]
WantedBy=multi-user.target
