# stage3_custom/02-wifi-connect/files/wifi-connect.service
[Unit]
Description=WiFi Captive Portal
After=NetworkManager.service
Before=nginx.service homebridge.service homebridge-config-ui-x.service homebridge-arch-check.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/wifi-connect-startup

[Install]
WantedBy=multi-user.target
